# hiwifi_http_filetransit/etc/init.d/lighttpd
#!/bin/sh /etc/rc.common

START=99

start() {
	lighttpd -f /etc/lighttpd/lighttpd.conf &
}

stop() {
	killall lighttpd
}

